# docs/pibo-cli/assets/services/notifier.service
[Unit]
Description=PIBo Docs Pull Notifier — HTTP trigger for immediate pull
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{NODE_BIN}} {{SCRIPTS_DIR}}/pibo-notifier.js
Restart=always
RestartSec=5
WorkingDirectory={{HOME}}
Environment=HOME={{HOME}}
Environment=PATH={{NODE_BIN_DIR}}:{{PATH}}
Environment=DOCS_PATH={{DOCS_PATH}}

[Install]
WantedBy=default.target
